# data/01Simple1.service
[Unit]
Description=
After=
Wants=

[Service]
User=root
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target