# data/Simple 2.service
[Unit]
Description=
After=
Wants=

[Service]
User=root
ExecStart=/bin/true
WorkingDirectory=/home

[Install]
WantedBy=multi-user.target